# docker/etc/wingmate/crontab
17 * * * * sleep 1
*/12 * * * * sleep 1
12,17,27 * * * * sleep 1